# scripts/com.tripline.daily-check.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<!--
  daily-check LaunchDaemon — 2026-05-11 從 LaunchAgent 升級
  目的：擺脫 user session XPC throttle（2026-05-07~05-10 連 4 天 launchd 沒 fire 推測為此）
  安裝/重新載入：bash scripts/install-daily-check-launchdaemon.sh
-->
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.tripline.daily-check</string>

    <key>UserName</key>
    <string>ray</string>

    <key>ProgramArguments</key>
    <array>
        <string>/bin/zsh</string>
        <string>/Users/ray/.local/bin/tripline-daily-check.sh</string>
    </array>

    <key>WorkingDirectory</key>
    <string>/Users/ray/Projects/trip-planner</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin:/Users/ray/.local/bin:/Users/ray/.bun/bin</string>
        <key>HOME</key>
        <string>/Users/ray</string>
        <key>USER</key>
        <string>ray</string>
        <key>LANG</key>
        <string>en_US.UTF-8</string>
    </dict>

    <key>StartCalendarInterval</key>
    <dict>
        <key>Hour</key>
        <integer>6</integer>
        <key>Minute</key>
        <integer>13</integer>
    </dict>

    <key>StandardOutPath</key>
    <string>/Users/ray/Projects/trip-planner/.context/daily-check-stdout.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/ray/Projects/trip-planner/.context/daily-check-stderr.log</string>
</dict>
</plist>
